# legacy-backup/com.willmclemore.projecttracker.wake.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.willmclemore.projecttracker.wake</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/willmclemore/wills_code/project_tracker/check-missed-jobs.sh</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>StartOnMount</key>
    <false/>
    <key>WakeFromSleep</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/willmclemore/wills_code/project_tracker/wake-check.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/willmclemore/wills_code/project_tracker/wake-check.log</string>
</dict>
</plist>
